# kickstart/etc/init/blink1.conf
description   "blink(1)"

start on (local-filesystems and net-device-up and runlevel [2345])

task
script
  /usr/local/bin/blink1-tool --rgb 00ff00
end script
