# kickstart/etc/gis.crontab
PATH=/usr/local/bin:/usr/sbin:/usr/bin:/bin

* * * * * apply-updates.sh
