# kickstart/etc/gis.crontab
PATH=/usr/local/bin:/usr/sbin:/usr/bin:/bin

* * * * * sleep 30 && apply-updates.sh
