# kickstart/etc/odm-web.upstart
description     "OpenDroneMap API - web"

start on (local-filesystems and started redis)
stop on runlevel [!2345]

respawn

script
  /usr/bin/docker start -a odm-web
end script

post-stop script
  /usr/bin/docker stop odm-web
end script
