# kickstart/etc/odm-worker.upstart
description     "OpenDroneMap API - worker"

start on (local-filesystems and started odm-web)
stop on runlevel [!2345]

respawn

script
  /usr/bin/docker start -a odm-worker
end script
