# kickstart/etc/omk-server.upstart
# omk-server
description     "OpenMapKitServer"

start on local-filesystems
stop on runlevel [!2345]

env PORT="{{omk_server_port}}"
env NODE_ENV="{{posm_env}}"

respawn

script
  /usr/bin/docker start -a omk-server
end script

post-stop script
  /usr/bin/docker stop omk-server
end script
