# kickstart/etc/redis.upstart
description     "Redis"

start on (local-filesystems and started docker)
stop on runlevel [!2345]

console log
respawn

script
  /usr/bin/docker start -a redis
end script

post-stop script
  /usr/bin/docker stop redis
end script
