# kickstart/etc/root.crontab
* * * * * /usr/bin/posm-imagery-updater > /opt/posm-www/posm_imagery.json
