# kickstart/etc/systemd/system/osmosis-replication.timer
[Unit]
Description=osmosis replication

[Timer]
OnBootSec=15
OnUnitActiveSec=60

[Install]
WantedBy=multi-user.target
