# kickstart/etc/systemd/system/wan-timeout.timer
[Unit]
Description=WAN DHCP timeout

[Timer]
OnActiveSec=15