# kickstart/etc/systemd/system/nodeodm.service.d/override.conf
[Service]
ExecStartPost=-/usr/local/bin/register-nodeodm-worker.sh
