# kickstart/etc/systemd/system/nodeodm.service.d/override.conf
[Unit]
After=docker.service network-online.target
Requires=docker.service network-online.target

[Service]
ExecStartPost=-/usr/local/bin/register-nodeodm-worker.sh
